# scripts/service/92_wss_parse_slot_duration.service
cd ..[Unit]
Description=WebSocket Slot Duration Parser (Group 1)
After=network.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=smilax
Group=smilax
WorkingDirectory=/home/smilax/trillium_api/data/monitoring/wss_slot_duration
ExecStart=/home/smilax/.python_env/bin/python /home/smilax/trillium_api/scripts/python/92_wss_parse_slot_duration.py --group 1
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

# Resource limits
LimitNOFILE=65536
TimeoutStartSec=30
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/smilax/log /home/smilax/trillium_api/data/monitoring/wss_slot_duration

[Install]
WantedBy=multi-user.target